# deploy/uci-defaults/99-portcullis
#!/bin/sh
# First-boot bootstrap for portcullis (TDD §10). Runs once on first boot, then
# OpenWrt deletes it. Idempotent: safe if the fleet pipeline re-runs it.
#
# What it sets up that the package install can't: the dedicated daemon user, the
# tmpfs state dir, the dnsmasq-full garden nftset wiring, and a placeholder
# secrets dir. Per-store identity (store_id, endpoint, HMAC key, mTLS client
# cert/key + pinned CP server CA) is injected by the fleet provisioning
# pipeline, NOT baked into the image (§13).

set -e

# 1) Dedicated unprivileged user/group (no shell, no home). The daemon drops to
#    this user with CAP_NET_ADMIN via procd (see capabilities/portcullis.json).
if ! grep -q '^portcullis:' /etc/passwd 2>/dev/null; then
	echo 'portcullis:x:err:err:portcullis:/var/run/portcullis:/bin/false' >/dev/null 2>&1 || true
	group_add portcullis 2>/dev/null || true
	user_add portcullis '' '' 'portcullis' /var/run/portcullis /bin/false 2>/dev/null || true
fi

# 2) Secrets dir (HMAC key + mTLS material land here, 0600, owned by the daemon).
#    The pipeline drops the real key/certs; we only prepare the directory.
mkdir -p /etc/portcullis/tls
chmod 0700 /etc/portcullis /etc/portcullis/tls
chown portcullis:portcullis /etc/portcullis /etc/portcullis/tls 2>/dev/null || true

# 3) tmpfs runtime state dir (never NAND — §5.4).
mkdir -p /tmp/portcullis

# 4) Walled garden: the ENGINE is the single owner of the dnsmasq garden conf at
#    runtime (portcullis-garden). It probes dnsmasq for ipset/nftset support before
#    writing (so a stock/slim dnsmasq is never handed a fatal `ipset=`/`nftset=`
#    directive), picks the set family that MATCHES the chosen firewall backend
#    (nft→nftset / ipset→ipset — a hardcoded seed here could disagree and silently
#    empty the garden), validates every FQDN, and writes atomically.
#    A pre-engine seed is therefore both unnecessary (no engine running = no gate,
#    so the garden is moot) and dangerous (unprobed + wrong-family + on flash where
#    the daemon's self-heal can't reach it). We deliberately DO NOT write one here.
#    See RC/audit #5 (unprobed fatal conf on fresh boot) and #6 (family mismatch).

# 5) Enable + start the service.
/etc/init.d/portcullis enable 2>/dev/null || true

exit 0
